# ansible/roles/influxdb-docker-backup/templates/influx-backup.service
[Unit]
Description=Backup influxdb
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/home/sascha/backup.sh
User=sascha
Group=sascha
